# scripts/launchd/com.mdeane.article-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mdeane.article-ingest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/matthewdeane/Documents/Data Science/python/_projects/____p-library/library-articles/scripts/ingest.py</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/matthewdeane/Library/Logs/article-ingest.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/matthewdeane/Library/Logs/article-ingest.err</string>

    <key>WorkingDirectory</key>
    <string>/Users/matthewdeane/Documents/Data Science/python/_projects/____p-library/library-articles</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
